# files/etc/uci-defaults/30-change_hostname
#!/bin/sh
[ "$(uci -q get system.@system[0].hostname)" = "FoxHole" ] && exit 0

uci -q batch <<-EOF >/dev/null
	set system.@system[0].hostname='FoxHole'
	commit system
EOF

exit 0