# files/etc/uci-defaults/31-change_banner
#!/bin/sh
if grep -q "Freedom for foxes!" /etc/banner >/dev/null; then
    exit 0
fi

# Очищаем файл banner
echo "" > /etc/banner

# Записываем новое содержимое в файл banner
cat << EOF > /etc/banner
    ▄████████  ▄██████▄  ▀████    ▐████▀    ▄█    █▄     ▄██████▄   ▄█          ▄████████ 
   ███    ███ ███    ███   ███▌   ████▀    ███    ███   ███    ███ ███         ███    ███ 
   ███    █▀  ███    ███    ███  ▐███      ███    ███   ███    ███ ███         ███    █▀  
  ▄███▄▄▄     ███    ███    ▀███▄███▀     ▄███▄▄▄▄███▄▄ ███    ███ ███        ▄███▄▄▄     
 ▀▀███▀▀▀     ███    ███    ████▀██▄     ▀▀███▀▀▀▀███▀  ███    ███ ███       ▀▀███▀▀▀     
   ███        ███    ███   ▐███  ▀███      ███    ███   ███    ███ ███         ███    █▄  
   ███        ███    ███  ▄███     ███▄    ███    ███   ███    ███ ███▌    ▄   ███    ███ 
   ███         ▀██████▀  ████       ███▄   ███    █▀     ▀██████▀  █████▄▄██   ██████████ 
                                                                   ▀                      

   Freedom for foxes!                                                          
 ---------------------------------------------------------------------------------------
 %D %V, %C
 ---------------------------------------------------------------------------------------
EOF

exit 0
