# files/etc/uci-defaults/40-change_ruantiblock_config
#!/bin/sh
[ "$(uci -q get ruantiblock.config.proxy_mode)" = "2" ] && exit 0

uci -q batch <<-EOF >/dev/null
	uci set ruantiblock.config.proxy_mode="2"
	commit ruantiblock
EOF

exit 0